# utility/scripts/pycon2017.conf
[program:pycon2017_{{ environment }}]
command = {{ run_root }}/gunicorn_run.sh                     ; Command to start app
user = deploy                                                ; User to run as
stdout_logfile = {{ logs_root }}/gunicorn_supervisor.log     ; Where to write log messages
redirect_stderr = true                                       ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8              ; Set UTF-8 as default encoding
